# krayin-prod/deploy/supervisor/krayin-queue.conf
[program:krayin-queue]
command=php /var/www/krayin/artisan queue:work --tries=3
user=www-data
numprocs=1
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/krayin/storage/logs/queue.log
